# cyberswarm_cli/knowledge/threat-intelligence.json
{
  "indicators": [
    {
      "type": "ip",
      "value": "192.168.1.100",
      "description": "Known C2 server communication endpoint",
      "severity": "High",
      "first_seen": "2024-01-01",
      "last_seen": "2024-06-15",
      "tags": ["malware", "c2", "cobalt-strike"],
      "mitre_techniques": ["T1071", "T1572"]
    },
    {
      "type": "ip",
      "value": "10.0.0.55",
      "description": "Internal lateral movement pivot point",
      "severity": "Critical",
      "first_seen": "2024-03-10",
      "last_seen": "2024-06-20",
      "tags": ["lateral-movement", "psexec", "wmiexec"],
      "mitre_techniques": ["T1021", "T1210"]
    },
    {
      "type": "ip",
      "value": "203.0.113.42",
      "description": "External scanning source associated with APT campaigns",
      "severity": "Medium",
      "first_seen": "2024-02-15",
      "last_seen": "2024-05-30",
      "tags": ["scanning", "reconnaissance"],
      "mitre_techniques": ["T1595", "T1046"]
    },
    {
      "type": "domain",
      "value": "update-service.example.net",
      "description": "Domain used for malware C2 disguised as software update service",
      "severity": "High",
      "first_seen": "2024-04-01",
      "last_seen": "2024-06-18",
      "tags": ["c2", "dns-tunnel", "dga"],
      "mitre_techniques": ["T1071", "T1572"]
    },
    {
      "type": "domain",
      "value": "cdn-assets.example.org",
      "description": "Phishing infrastructure hosting credential harvesting pages",
      "severity": "High",
      "first_seen": "2024-05-10",
      "last_seen": "2024-06-22",
      "tags": ["phishing", "credential-harvesting"],
      "mitre_techniques": ["T1566", "T1078"]
    },
    {
      "type": "hash",
      "value": "a1b2c3d4e5f6a7b8c9d0e1f2a3b4c5d6e7f8a9b0",
      "description": "SHA1 hash of Cobalt Strike Beacon payload",
      "severity": "Critical",
      "first_seen": "2024-03-20",
      "last_seen": "2024-06-15",
      "tags": ["cobalt-strike", "beacon", "malware"],
      "mitre_techniques": ["T1059", "T1055", "T1071"]
    },
    {
      "type": "hash",
      "value": "e5f6a7b8c9d0e1f2a3b4c5d6e7f8a9b0c1d2e3f4",
      "description": "SHA1 hash of custom credential dumping tool",
      "severity": "Critical",
      "first_seen": "2024-04-15",
      "last_seen": "2024-06-10",
      "tags": ["credential-dumping", "mimikatz-variant"],
      "mitre_techniques": ["T1003", "T1003.001"]
    },
    {
      "type": "url",
      "value": "https://update-service.example.net/api/v2/check",
      "description": "C2 beacon check-in URL",
      "severity": "High",
      "first_seen": "2024-04-05",
      "last_seen": "2024-06-18",
      "tags": ["c2", "beacon", "https"],
      "mitre_techniques": ["T1071.001"]
    },
    {
      "type": "file_path",
      "value": "C:\\Windows\\Temp\\svchost_update.exe",
      "description": "Malware masquerading as legitimate Windows process",
      "severity": "Critical",
      "first_seen": "2024-03-25",
      "last_seen": "2024-06-12",
      "tags": ["masquerading", "persistence", "dropper"],
      "mitre_techniques": ["T1036", "T1547"]
    },
    {
      "type": "email",
      "value": "admin-support@cdn-assets.example.org",
      "description": "Phishing email sender address used in spear-phishing campaigns",
      "severity": "Medium",
      "first_seen": "2024-05-10",
      "last_seen": "2024-06-22",
      "tags": ["phishing", "spear-phishing"],
      "mitre_techniques": ["T1566.001"]
    }
  ],
  "campaigns": [
    {
      "name": "APT-SIM-001",
      "description": "Simulated APT campaign targeting enterprise networks with multi-stage attack chain",
      "tactics": ["reconnaissance", "initial-access", "persistence", "lateral-movement", "exfiltration"],
      "techniques": ["T1595", "T1190", "T1053", "T1021", "T1041"],
      "severity": "High",
      "start_date": "2024-01-15",
      "status": "active",
      "targets": ["enterprise-networks", "web-applications"],
      "tools_used": ["nmap", "nuclei", "metasploit", "cobalt-strike", "impacket"]
    },
    {
      "name": "RANSOMWARE-SIM-002",
      "description": "Simulated ransomware operation using double extortion methodology",
      "tactics": ["initial-access", "execution", "persistence", "credential-access", "lateral-movement", "collection", "exfiltration", "impact"],
      "techniques": ["T1566", "T1059", "T1547", "T1003", "T1021", "T1119", "T1041", "T1486"],
      "severity": "Critical",
      "start_date": "2024-03-01",
      "status": "active",
      "targets": ["domain-controllers", "file-servers", "backup-systems"],
      "tools_used": ["cobalt-strike", "mimikatz", "bloodhound", "impacket"]
    },
    {
      "name": "SUPPLY-CHAIN-SIM-003",
      "description": "Simulated supply chain compromise targeting software update mechanisms",
      "tactics": ["initial-access", "execution", "persistence", "defense-evasion"],
      "techniques": ["T1195", "T1059", "T1547", "T1027"],
      "severity": "Critical",
      "start_date": "2024-04-01",
      "status": "monitoring",
      "targets": ["ci-cd-pipelines", "package-repositories", "update-servers"],
      "tools_used": ["nuclei", "sigma", "yara"]
    },
    {
      "name": "CREDENTIAL-HARVEST-004",
      "description": "Credential harvesting campaign using LLMNR/NBT-NS poisoning and Kerberoasting",
      "tactics": ["credential-access", "lateral-movement", "privilege-escalation"],
      "techniques": ["T1557", "T1558.003", "T1110.002", "T1003", "T1550"],
      "severity": "High",
      "start_date": "2024-05-01",
      "status": "active",
      "targets": ["active-directory", "kerberos", "ntlm"],
      "tools_used": ["responder", "impacket", "mimikatz", "hashcat", "bloodhound"]
    },
    {
      "name": "CLOUD-PIVOT-005",
      "description": "Simulated cloud infrastructure attack pivoting from on-prem compromise",
      "tactics": ["discovery", "lateral-movement", "collection", "exfiltration"],
      "techniques": ["T1087", "T1069", "T1021", "T1530", "T1567"],
      "severity": "High",
      "start_date": "2024-06-01",
      "status": "planned",
      "targets": ["azure-ad", "aws-iam", "cloud-storage"],
      "tools_used": ["bloodhound", "impacket", "metasploit"]
    }
  ],
  "actors": [
    {
      "name": "PHANTOM-BEAR",
      "description": "Advanced persistent threat group with state-sponsored capabilities targeting critical infrastructure",
      "aliases": ["APT-PHANTOM", "Dark Ursus", "Ghost Grizzly"],
      "sophistication": "High",
      "motivation": "Espionage",
      "origin": "Unknown",
      "ttps": ["spear-phishing", "zero-day-exploitation", "custom-malware", "living-off-the-land", "data-exfiltration"],
      "mitre_techniques": ["T1566", "T1190", "T1059", "T1055", "T1003", "T1071"],
      "target_sectors": ["energy", "government", "defense"],
      "tools_observed": ["cobalt-strike", "mimikatz", "custom-implants"],
      "confidence": 0.85
    },
    {
      "name": "CYBER-SYNDICATE-7",
      "description": "Financially motivated cybercriminal group specializing in ransomware and double extortion",
      "aliases": ["RansomCrew7", "DarkVault"],
      "sophistication": "Medium",
      "motivation": "Financial",
      "origin": "Unknown",
      "ttps": ["phishing", "ransomware", "data-exfiltration", "credential-stealing", "supply-chain"],
      "mitre_techniques": ["T1566", "T1486", "T1041", "T1003", "T1195"],
      "target_sectors": ["healthcare", "finance", "manufacturing"],
      "tools_observed": ["cobalt-strike", "metasploit", "mimikatz", "bloodhound"],
      "confidence": 0.78
    },
    {
      "name": "SILENT-SPIDER",
      "description": "Stealthy threat actor focused on long-term persistence and intelligence collection",
      "aliases": ["WebWeaver", "SilkThread"],
      "sophistication": "High",
      "motivation": "Espionage",
      "origin": "Unknown",
      "ttps": ["watering-hole", "custom-backdoors", "dns-tunneling", "process-injection", "timestomping"],
      "mitre_techniques": ["T1189", "T1059", "T1572", "T1055", "T1070"],
      "target_sectors": ["technology", "telecommunications", "research"],
      "tools_observed": ["sliver", "custom-implants", "impacket"],
      "confidence": 0.72
    },
    {
      "name": "VOLT-JACKAL",
      "description": "Opportunistic threat actor targeting exposed services and known vulnerabilities",
      "aliases": ["QuickStrike", "OpportunityKnocks"],
      "sophistication": "Low",
      "motivation": "Financial",
      "origin": "Unknown",
      "ttps": ["vulnerability-scanning", "exploit-kits", "cryptomining", "credential-stuffing"],
      "mitre_techniques": ["T1595.002", "T1190", "T1110", "T1496"],
      "target_sectors": ["small-business", "education", "any-exposed"],
      "tools_observed": ["nuclei", "masscan", "metasploit"],
      "confidence": 0.90
    }
  ],
  "detection_rules": [
    {
      "rule_id": "SIGMA-001",
      "name": "Suspicious LSASS Access",
      "description": "Detects potential credential dumping via LSASS memory access",
      "format": "sigma",
      "mitre_technique": "T1003.001",
      "severity": "Critical",
      "data_source": "Sysmon Event ID 10",
      "false_positive_rate": "low"
    },
    {
      "rule_id": "SIGMA-002",
      "name": "Suspicious Scheduled Task Creation",
      "description": "Detects creation of scheduled tasks commonly used for persistence",
      "format": "sigma",
      "mitre_technique": "T1053",
      "severity": "High",
      "data_source": "Windows Security Event ID 4698",
      "false_positive_rate": "medium"
    },
    {
      "rule_id": "SIGMA-003",
      "name": "DNS Tunneling Detection",
      "description": "Detects potential DNS tunneling based on query length and frequency anomalies",
      "format": "sigma",
      "mitre_technique": "T1572",
      "severity": "High",
      "data_source": "DNS query logs",
      "false_positive_rate": "medium"
    },
    {
      "rule_id": "YARA-001",
      "name": "CobaltStrike Beacon Detection",
      "description": "Detects Cobalt Strike Beacon payloads in memory or on disk",
      "format": "yara",
      "mitre_technique": "T1071",
      "severity": "Critical",
      "data_source": "File/memory scanning",
      "false_positive_rate": "low"
    },
    {
      "rule_id": "YARA-002",
      "name": "Mimikatz String Detection",
      "description": "Detects Mimikatz tool presence via string patterns",
      "format": "yara",
      "mitre_technique": "T1003",
      "severity": "Critical",
      "data_source": "File/memory scanning",
      "false_positive_rate": "low"
    },
    {
      "rule_id": "SURICATA-001",
      "name": "Metasploit Reverse Shell",
      "description": "Detects Metasploit reverse shell network signatures",
      "format": "suricata",
      "mitre_technique": "T1059",
      "severity": "Critical",
      "data_source": "Network traffic",
      "false_positive_rate": "low"
    },
    {
      "rule_id": "SIGMA-004",
      "name": "Kerberoasting Activity",
      "description": "Detects Kerberoasting by monitoring for anomalous TGS requests",
      "format": "sigma",
      "mitre_technique": "T1558.003",
      "severity": "High",
      "data_source": "Windows Security Event ID 4769",
      "false_positive_rate": "low"
    },
    {
      "rule_id": "SIGMA-005",
      "name": "Process Injection via CreateRemoteThread",
      "description": "Detects process injection attempts using CreateRemoteThread API",
      "format": "sigma",
      "mitre_technique": "T1055",
      "severity": "Critical",
      "data_source": "Sysmon Event ID 8",
      "false_positive_rate": "medium"
    }
  ],
  "threat_feeds": [
    {
      "name": "MISP Default Feeds",
      "type": "ioc_feed",
      "format": "misp-json",
      "update_frequency": "hourly",
      "categories": ["ip", "domain", "hash", "url"]
    },
    {
      "name": "Sigma Community Rules",
      "type": "detection_rules",
      "format": "sigma-yaml",
      "update_frequency": "daily",
      "categories": ["detection", "windows", "linux", "network"]
    },
    {
      "name": "YARA Community Rules",
      "type": "detection_rules",
      "format": "yara",
      "update_frequency": "weekly",
      "categories": ["malware", "exploits", "webshells"]
    },
    {
      "name": "Emerging Threats Open",
      "type": "ids_rules",
      "format": "suricata",
      "update_frequency": "daily",
      "categories": ["network", "malware", "exploit"]
    },
    {
      "name": "MITRE ATT&CK STIX Data",
      "type": "framework",
      "format": "stix-2.1",
      "update_frequency": "quarterly",
      "categories": ["techniques", "tactics", "mitigations", "groups"]
    }
  ],
  "mitre_coverage_baseline": {
    "Reconnaissance": 0.45,
    "Resource Development": 0.20,
    "Initial Access": 0.65,
    "Execution": 0.70,
    "Persistence": 0.60,
    "Privilege Escalation": 0.55,
    "Defense Evasion": 0.40,
    "Credential Access": 0.75,
    "Discovery": 0.50,
    "Lateral Movement": 0.60,
    "Collection": 0.35,
    "Command and Control": 0.55,
    "Exfiltration": 0.30,
    "Impact": 0.45
  }
}
